# deploy/stockbot.service
# systemd 服務：跑一次掃描就退出 (Type=oneshot)，由 stockbot.timer 定時觸發。
# 安裝：
#   sudo cp deploy/stockbot.service /etc/systemd/system/
#   依實際路徑/使用者修改下面的 User、WorkingDirectory、ExecStart
[Unit]
Description=Taiwan stock famous-strategy scan (run once)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/stock
# 實單金鑰請放在這個 EnvironmentFile（chmod 600，且已被 .gitignore 忽略）
EnvironmentFile=-/home/YOUR_USER/stock/.env
# 四種模式（擇一，把對應那行的 # 拿掉）。--regime=大盤空頭禁買 (回測證實能保本，建議保留)：
# --universe tw50 = 餵 50 大權值股當候選池，系統自動挑「當下合格且訊號最強」的
# --max-positions 5 = 最多同時持有 5 檔 (5萬建議搭 --budget 10000，每檔約1萬)
#  A. 純通知 dry-run（不下單，只推 Telegram）—— 最安全
# ExecStart=/home/YOUR_USER/stock/.venv/bin/python main.py scan --strategy lynch --source finmind --universe tw50 --realtime --regime --max-positions 5 --budget 10000 --notify
#  P. 本地持久化模擬盤（假錢、自己記帳、收盤後用當天日 K 撮合）—— 目前選用，最推薦的空跑方式。
#     持倉存 paper_account.json、跨執行累積、數字乾淨可信 (不用永豐模擬盤——其回報實測不可靠)。
#     搭 timer 收盤後(14:00)一天跑一次；lynch 是基本面策略不需盤中即時價，也省 FinMind 額度。
#     --cash 30000 只在首次建檔時生效；要重置就刪掉 paper_account.json。
#     雙策略配置 (2026-07-02 回測+walkforward 結論)：lynch 防守核心 3 萬/3 檔，
#     livermore 進攻衛星 2 萬/2 檔獨立帳戶 → 見 stockbot-livermore.service/.timer。
ExecStart=/home/YOUR_USER/stock/.venv/bin/python main.py scan --strategy lynch --source finmind --universe tw50 --regime --max-positions 3 --budget 10000 --cash 30000 --paper --notify
#  B. Shioaji 模擬盤全自動（假錢自動下單）—— 不建議：永豐模擬盤持倉/成交回報實測會亂噴。
# ExecStart=/home/YOUR_USER/stock/.venv/bin/python main.py scan --strategy lynch --source finmind --universe tw50 --realtime --regime --max-positions 5 --budget 10000 --live --notify
#  C. 實單全自動（真錢！需先設定憑證 SHIOAJI_CA_PATH 等）—— 確認 P 模式空跑穩定才開
# ExecStart=/home/YOUR_USER/stock/.venv/bin/python main.py scan --strategy lynch --source finmind --universe tw50 --realtime --regime --max-positions 5 --budget 10000 --live --real-account --notify

# 資源上限：避免萬一吃爆記憶體影響同台的其他服務（如 pionex-bot）
MemoryMax=400M
CPUQuota=50%
Nice=10

[Install]
WantedBy=multi-user.target
