# deploy/stockbot.service
# systemd 服務：跑一次掃描就退出 (Type=oneshot)，由 stockbot.timer 定時觸發。
# 安裝：
#   sudo cp deploy/stockbot.service /etc/systemd/system/
#   依實際路徑/使用者修改下面的 User、WorkingDirectory、ExecStart
[Unit]
Description=Taiwan stock famous-strategy scan (run once)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/stock
# 實單金鑰請放在這個 EnvironmentFile（chmod 600，且已被 .gitignore 忽略）
EnvironmentFile=-/home/YOUR_USER/stock/.env
# 三種模式（擇一，把對應那行的 # 拿掉）。--regime=大盤空頭禁買 (回測證實能保本，建議保留)：
#  A. 純通知 dry-run（不下單，只推 Telegram）—— 最安全
# ExecStart=/home/YOUR_USER/stock/.venv/bin/python main.py scan --strategy lynch --source finmind --realtime --regime --notify
#  B. Shioaji 模擬盤全自動（假錢自動下單）—— 目前選用，--live 但無 --real-account=模擬
ExecStart=/home/YOUR_USER/stock/.venv/bin/python main.py scan --strategy lynch --source finmind --realtime --regime --live --notify
#  C. 實單全自動（真錢！需先設定憑證 SHIOAJI_CA_PATH 等）—— 確認模擬盤穩定才開
# ExecStart=/home/YOUR_USER/stock/.venv/bin/python main.py scan --strategy lynch --source finmind --realtime --regime --live --real-account --notify

# 資源上限：避免萬一吃爆記憶體影響同台的其他服務（如 pionex-bot）
MemoryMax=400M
CPUQuota=50%
Nice=10

[Install]
WantedBy=multi-user.target
