# deploy/stockbot-listen.service
# systemd 服務：常駐監聽 Telegram 指令，動態調整 budget / max_positions / pause。
# 安裝：
#   sed -i "s/YOUR_USER/$(whoami)/g" deploy/stockbot-listen.service
#   sudo cp deploy/stockbot-listen.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now stockbot-listen.service
#   journalctl -u stockbot-listen.service -f      # 看即時 log
[Unit]
Description=Taiwan stock bot - Telegram command listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER/stock
EnvironmentFile=-/home/YOUR_USER/stock/.env
# --paper：/holdings /sell 對本地持久化模擬盤帳戶。預設同時看兩個策略帳戶
# (lynch=paper_account.json + livermore=paper_livermore.json)：/holdings 合併顯示
# 各帳戶持倉+現金+總資產；/sell 自動找到持有那檔的帳戶賣。還沒建檔的帳戶不會顯示。
# 若 scan 用 Shioaji 模擬盤/實單 (B/C 模式)，就把 --paper 拿掉，改對券商帳戶。
ExecStart=/home/YOUR_USER/stock/.venv/bin/python main.py listen --paper
Restart=always
RestartSec=10
MemoryMax=200M

[Install]
WantedBy=multi-user.target
