# deploy/stockbot.timer
# systemd 計時器：台股收盤後（14:00 後）每日跑一次掃描。
# 安裝：
#   sudo cp deploy/stockbot.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now stockbot.timer
#   systemctl list-timers stockbot*      # 確認下次觸發時間
[Unit]
Description=Run stock scan daily after market close

# 注意：請先把 VM 時區設成台北，否則時間會對不上：
#   sudo timedatectl set-timezone Asia/Taipei
[Timer]
# 台股盤中 (09:00–13:30) 每 5 分鐘掃一次。
# 09..13:00/5 = 09~13 時、每 5 分鐘 (含 13:00~13:30，13:30 後幾筆無妨)。
OnCalendar=Mon-Fri 09..13:00/5
Persistent=false

[Install]
WantedBy=timers.target
